# service/systemd/btrfs-snapshot.timer
[Unit]
Description=Run Btrfs snapshot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

