# service/systemd/btrfs-snapshot.service
[Unit]
Description=Btrfs snapshot everyday

[Service]
Type=oneshot
ExecStart=/usr/bin/zsh /opt/btrfs-snapshot

